# deploy/site1vpn-restart.path
[Unit]
Description=Watch for 1VPN restart requests

[Path]
PathExists=/var/www/site1vpn/runtime/restart.request
Unit=site1vpn-restart.service

[Install]
WantedBy=multi-user.target
