# deploy/site1vpn.service
[Unit]
Description=1VPN Next.js application
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=site1vpn
Group=site1vpn
WorkingDirectory=/var/www/site1vpn
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=NODE_OPTIONS=--max-old-space-size=768
ExecStart=/usr/bin/node /var/www/site1vpn/node_modules/next/dist/bin/next start --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=20
LimitNOFILE=65535
NoNewPrivileges=true
UMask=0027

[Install]
WantedBy=multi-user.target
