# deploy/site1vpn.service
[Unit]
Description=1VPN Next.js application
After=network.target

[Service]
Type=simple
User=site1vpn
Group=site1vpn
WorkingDirectory=/var/www/site1vpn
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/www/site1vpn/.next /var/www/site1vpn/prisma

[Install]
WantedBy=multi-user.target
