# deploy/site1vpn-healthcheck.service
[Unit]
Description=Check and recover 1VPN application and HTTPS
After=network-online.target site1vpn.service nginx.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /var/www/site1vpn/deploy/site1vpn-healthcheck.sh
TimeoutStartSec=75
