# deploy/site1vpn-recover.service
[Unit]
Description=Restart 1VPN after a failed health check

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart site1vpn.service
